# stage2/02-net-tweaks/files/ScanAndConnect.service
[Unit]
Description=Scan and Connect to 100- wifi
Wants=network.target
Before=dhcpcd5.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/captivePortal/WifiScanner
EnvironmentFile=/opt/captivePortal/WifiScanner/ScanAndConnect.conf
ExecStart=/opt/captivePortal/WifiScanner/ScanAndConnect.sh $ESSID $ESSID_PASSWORD
Restart=on-abort

[Install]
WantedBy=multi-user.target
